# scoreo2.service.in
[Unit]
Description=@PROJECT_NAME@
After=network.target

[Service]
Environment=LD_LIBRARY_PATH=@OPT_PREFIX@/lib
ExecStart=@OPT_PREFIX@/bin/starter
Restart=on-failure
User=root
WorkingDirectory=@OPT_PREFIX@

[Install]
WantedBy=multi-user.target
